# ansible/files/neurosimo-logs.service.j2
[Unit]
Description=neurosimo log follower
After=neurosimo.service docker.service
Requires=docker.service
PartOf=neurosimo.service

StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User={{ target_user }}
WorkingDirectory={{ repo_root }}

Environment=COMPOSE_PROJECT_NAME=neurosimo

# Stream logs to journald. Use the same profile selection as runtime.
{% if use_gpu | bool %}
ExecStart=/usr/bin/docker compose -f {{ repo_root }}/docker-compose.yml --profile gpu logs -f --tail=200
{% else %}
ExecStart=/usr/bin/docker compose -f {{ repo_root }}/docker-compose.yml --profile non-gpu logs -f --tail=200
{% endif %}

# Terminate immediately on stop since logs command doesn't handle SIGTERM gracefully
TimeoutStopSec=2
KillMode=mixed

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
